# notebooklm_scrape/com.notebooklm.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notebooklm.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/z/Desktop/PersonalProjects/ClaudeProjects/notebooklm_scrape/daily_sync.sh</string>
    </array>

    <!-- Run daily at 6 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/z/Desktop/PersonalProjects/ClaudeProjects/notebooklm_scrape/logs/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/z/Desktop/PersonalProjects/ClaudeProjects/notebooklm_scrape/logs/sync_error.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
